# systemd/dola-fetch-service.service
[Unit]
Description=Dola fetch task service
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/dola-fetch-service
Environment=DOLA_DATA_DIR=/var/lib/dola-fetch-service
Environment=DOLA_CONFIG_PATH=/var/lib/dola-fetch-service/config.json
ExecStart=/opt/dola-fetch-service/.venv/bin/python /opt/dola-fetch-service/run.py
Restart=always
RestartSec=3
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
